# update_pi/auto-update.service
[Unit]
Description=Auto Update Host Service

[Service]
Type=simple
ExecStart=/PATH_TO_UPDATE_SCRIPT/
User=root
