# update_pi/auto-update.service
[Unit]
Description=Auto Update Host Service

[Service]
Type=simple
ExecStart=/home/marius/update.sh
User=root
